# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:bottel]
command=python /app/bottel.py

[program:telebot]
command=python /app/telebot.py
